# cronjob-sa-update.service
[Unit]
Description=Spamassassin Rules Updates

[Service]
ExecStart=/usr/share/spamassassin/sa-update.cron
