# cronjob-sa-update.timer
[Unit]
Description=Spamassassin Rules Updates timer

[Timer]
OnCalendar=4:10

[Install]
WantedBy=cronjobs.target
